# supervisord.conf
[supervisord]
nodaemon=true

[program:git]
command=/scripts/run

[program:http]
command=bundle exec puma scripts/config.ru